# systemd/lcctop-codex.service
[Unit]
Description=lcctop Codex session watcher
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/lcctop-codex
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
